# rk3128/android_source/m3filesystem/inittab
#如果内核没有选中->device drivers->Character devices->Virtual terminal 
#->Support for console on virtual terminal 
#则需要指定console为终端设备

::sysinit:/etc/init.d/rcS     
console::askfirst:-/bin/sh
::ctrlaltdel:-/sbin/reboot
::shutdown:/sbin/swapoff -a
::shutdown:/bin/umount -a -r
::restart:/sbin/init
